# ansible/roles/bootstrap_k3s/templates/k3s.service.j2
{{ ansible_managed | comment }}
# Mirrors upstream github.com/k3s-io/k3s/blob/master/k3s.service. Type=notify,
# KillMode=process, Delegate=yes, TimeoutStartSec=0, Restart=always — all carry
# over unchanged. Substrate-required flags live OUTSIDE the public defaults:
#   * `--disable-cloud-controller` — hardcoded in ExecStart (plan §2.4).
#   * `--kubelet-arg=feature-gates=KubeletInUserNamespace=true` — hardcoded
#     in ExecStart (plan §2.8 unprivileged LXC).
#   * `--disable=traefik` + `--disable=servicelb` — sourced from
#     `_bootstrap_k3s_required_disable_components` in vars/main.yml
#     (plan §2.9 / §5.5 deliver ingress + LoadBalancer via Terraform Helm).
#   * `--cluster-cidr` / `--service-cidr` / `--cluster-dns` —
#     dual-stack values sourced from `vars/main.yml` (plan §13.9 — pod
#     CNI participates in dual-stack so CAPI controllers can reach
#     IPv6-addressed workload-cluster LB endpoints CAPN derives).
#   * `--flannel-ipv6-masq` — required for pod-originated IPv6 SNAT.
# See the role README + vars/main.yml for the full rationale.

[Unit]
Description=Lightweight Kubernetes (k3s-server)
Documentation=https://k3s.io
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
EnvironmentFile=-{{ bootstrap_k3s_env_path }}
KillMode=process
Delegate=yes
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
TimeoutStartSec=0
Restart=always
RestartSec=5s
ExecStart={{ _bootstrap_k3s_required_launcher_path }} \
    --write-kubeconfig-mode={{ bootstrap_k3s_write_kubeconfig_mode }} \
    --disable-cloud-controller \
    --kubelet-arg=feature-gates={{ (['KubeletInUserNamespace=true'] + bootstrap_k3s_extra_kubelet_feature_gates) | join(',') }} \
    --cluster-cidr={{ _bootstrap_k3s_required_cluster_cidr_v4 }},{{ _bootstrap_k3s_required_cluster_cidr_v6 }} \
    --service-cidr={{ _bootstrap_k3s_required_service_cidr_v4 }},{{ _bootstrap_k3s_required_service_cidr_v6 }} \
    --cluster-dns={{ _bootstrap_k3s_required_cluster_dns_v4 }},{{ _bootstrap_k3s_required_cluster_dns_v6 }} \
{% if _bootstrap_k3s_required_flannel_ipv6_masq %}
    --flannel-ipv6-masq \
{% endif %}
{% if _bootstrap_k3s_required_disable_network_policy %}
    --disable-network-policy \
{% endif %}
{% for component in (_bootstrap_k3s_required_disable_components + bootstrap_k3s_extra_disable_components) %}
    --disable={{ component }} \
{% endfor %}
{% for san in bootstrap_k3s_tls_san %}
    --tls-san={{ san }} \
{% endfor %}
{% if bootstrap_k3s_token | length > 0 %}
    --token={{ bootstrap_k3s_token }} \
{% endif %}
{% for arg in bootstrap_k3s_extra_args %}
    {{ arg }} \
{% endfor %}

[Install]
WantedBy=multi-user.target
